# files/.config/systemd/user/mako.service
[Unit]
Description=A lightweight Wayland notification daemon
Documentation=man:mako(1)
PartOf=sway-session.target

[Service]
Type=dbus
ExecStart=/usr/bin/mako
BusName=org.freedesktop.Notifications

[Install]
WantedBy=sway-session.target

