# home/dot_config/systemd/user/clean-teams-pwa.service
[Unit]
Description=Cleanup of Teams PWA cache and IndexedDB

[Service]
Type=oneshot
ExecStart=%h/.local/bin/clean-teams-pwa
